# services/certbot/auto-renew.cron
#!/bin/bash
cd $(dirname $0)
sleep $(awk 'BEGIN{srand(); print int(rand()*(3600+1))}')
exec docker compose run --rm -it certbot renew -q
